# deploy/systemd/paperclip.service
[Unit]
Description=Paperclip orchestration server
After=network-online.target openclaw.service
Wants=network-online.target

[Service]
Type=simple
User=paperclip
Group=paperclip
WorkingDirectory=/opt/openclaw-stack/paperclip/current
ExecStart=/usr/bin/node server/dist/index.js
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30
Environment=NODE_ENV=production
Environment=HOME=/var/lib/paperclip
Environment=PATH=/usr/bin:/usr/local/bin:/bin
Environment=PAPERCLIP_PORT=3100

[Install]
WantedBy=multi-user.target
